# launchd/com.stocktracking.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stocktracking.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Volumes/外接硬碟/stock-tracking/run.js</string>
  </array>

  <!-- 每天 14:00 執行 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>14</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Volumes/外接硬碟/stock-tracking</string>

  <key>StandardOutPath</key>
  <string>/Volumes/外接硬碟/stock-tracking/stock-tracking.log</string>
  <key>StandardErrorPath</key>
  <string>/Volumes/外接硬碟/stock-tracking/stock-tracking.log</string>

  <!-- 環境變數：確保 PATH 找得到 node -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
